# deploy/secretary-sync.timer
[Unit]
Description=Run the OpenCouncil secretary sync every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Unit=secretary-sync.service

[Install]
WantedBy=timers.target
